# Library/LaunchAgents/com.msrl.dotfiles.clean-downloads.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.msrl.dotfiles.clean-downloads</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shields/bin/clean_downloads.py</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>27</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/shields/Library/Logs/com.msrl.dotfiles.clean-downloads.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/shields/Library/Logs/com.msrl.dotfiles.clean-downloads.err</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
